# server/task-tracker.service
# task-tracker.service
# systemd unit for the Avalon Task Tracker server on EC2.
#
# Install:
#   sudo cp task-tracker.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable task-tracker
#   sudo systemctl start task-tracker

[Unit]
Description=Avalon Task Tracker Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/task-tracker

Environment="EC2_RELAY_URL=http://127.0.0.1:8766/send-signal"
Environment="RELAY_TOKEN=CHANGE_ME_shared_secret"
Environment="STATE_FILE=/home/ubuntu/task-tracker/state.json"

ExecStart=/home/ubuntu/task-tracker/venv/bin/python server.py

Restart=always
RestartSec=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=task-tracker

[Install]
WantedBy=multi-user.target
