# server/task-tracker.service
# task-tracker.service
# systemd unit for the Avalon Task Tracker Pi server.
#
# Install:
#   sudo cp task-tracker.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable task-tracker
#   sudo systemctl start task-tracker

[Unit]
Description=Avalon Task Tracker Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/task-tracker/server

# TODO: replace CHANGE_ME values with real secrets, or use an EnvironmentFile
Environment="EC2_RELAY_URL=http://13.58.219.0:8766/send-signal"
Environment="RELAY_TOKEN=CHANGE_ME_shared_secret"
Environment="STATE_FILE=/var/lib/task-tracker/state.json"

ExecStart=/home/pi/task-tracker/venv/bin/python server.py

Restart=always
RestartSec=5

# Ensure state dir exists with correct ownership
ExecStartPre=/bin/mkdir -p /var/lib/task-tracker
ExecStartPre=/bin/chown pi:pi /var/lib/task-tracker

StandardOutput=journal
StandardError=journal
SyslogIdentifier=task-tracker

[Install]
WantedBy=multi-user.target
